# etc/systemd/system/getty@tty7.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin wx672 --noclear %I $TERM
