# etc/systemd/system/android-backup@.service
[Service]
Type=oneshot
ExecStart=/home/wx672/.bin/android-backup.sh %I
