# dot.config/systemd/user/sandbox/default.target
[Unit]
Description=user services
Wants=emacs.service 
AllowIsolate=yes
