# dot.config/systemd/user/sandbox/post-wm.target
[Unit]
Description=stuffs after wm (chrome)
After=wm.target
AllowIsolate=yes

[Install]
Alias=default.target

