# dot.config/systemd/user/sandbox/pre-wm.target
[Unit]
Description=stuffs before wm (fcitx)
Before=wm.target
AllowIsolate=yes

[Install]
Alias=default.target
