# dot.config/systemd/user/sandbox/wm.target
[Unit]
Description=window manager
Requires=sawfish.service
Before=tools.target
AllowIsolate=yes

[Install]
Alias=default.target
